# _ops/supervisord.conf
[program:gunicorn-snake]
ommand = /home/ubuntu/anaconda3/envs/snake/bin/gunicorn snake.wsgi:application -w 4 -b unix:/tmp/gunicorn_snake.tmp --chdir=/home/ubuntu/prjs/snake/djangohome
